# cbuscgate2mqtt/rootfs/etc/s6-overlay/s6-rc.d/cgate/run
#!/usr/bin/with-contenv bashio
### #!/command/execlineb -P
bashio::log.info "Starting cgate"
cd /usr/local/bin/cgate
/usr/bin/java -Djava.awt.headless=true -jar -noverify /usr/local/bin/cgate/cgate.jar
