# cbuscgate2mqtt/rootfs/etc/s6-overlay/s6-rc.d/cgate/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Preparing to start C-Gate service..."

set -e

# Get config values into variables
CBUSNAME=$(bashio::config 'cbusname')
USEEVENTFILE=$(bashio::config 'useeventfile')
EVENTLEVEL=$(bashio::config 'eventlevel')

if ! [ -f /share/cgate/tag/${CBUSNAME}.xml ]
then
  bashio::log.fatal "Cbus project file /share/cgate/tag/${CBUSNAME}.xml not found!"
fi 

# Create symlinks for project and log directories
rm -r /usr/local/bin/cgate/tag 
ln -s /share/cgate/tag /usr/local/bin/cgate/tag
mkdir -p /share/cgate/logs
ln -s /share/cgate/logs /usr/local/bin/cgate/logs

bashio::log.info "Customising cgate C-GateConfig.txt"
sed -i "/project.default=/c\project.default=${CBUSNAME}" /usr/local/bin/cgate/config/C-GateConfig.txt
sed -i "/project.start=/c\project.start=${CBUSNAME}" /usr/local/bin/cgate/config/C-GateConfig.txt
if [ $USEEVENTFILE == true ]
then
    sed -i "/use-event-file=/c\use-event-file=yes" /usr/local/bin/cgate/config/C-GateConfig.txt
else
    sed -i "/use-event-file=/c\use-event-file=no" /usr/local/bin/cgate/config/C-GateConfig.txt
fi
sed -i "/event-file.event-level=/c\event-file.event-level=${EVENTLEVEL}" /usr/local/bin/cgate/config/C-GateConfig.txt

cd /usr/local/bin/cgate

bashio::log.info "Starting C-Gate!"
/usr/bin/java -Djava.awt.headless=true -jar -noverify /usr/local/bin/cgate/cgate.jar
